# dot_config/systemd/user/drive2pass-sync.service
[Unit]
Description=Sync Cloud Passwords Database to Local cache

[Service]
Type=oneshot
ExecStart=rclone sync -P Drive:/Passwords.kdbx %h
