# home/dot_config/systemd/user/goimapnotify@university.service.d/override.conf
[Service]
Environment="NOTMUCH_CONFIG=%E/notmuch/notmuchrc"
Environment="PASSWORD_STORE_DIR=%h/.local/share/pass"