# home/dot_config/systemd/user/pass2drive-sync.service
[Unit]
Description=Sync Local Passwords Database to the Cloud

[Service]
Type=oneshot
ExecStart=rclone sync -P %h/Passwords.kdbx Drive:
