# home/dot_config/systemd/user/pass2drive-sync.path
[Unit]
Description=Monitor ~/Passwords.kdbx for changes and sync to cloud

[Path]
PathModified=%h/Passwords.kdbx

[Install]
WantedBy=default.target
